# etc/upstart/noodle.conf
# noodle - noodle job file
description "noodle"
author "Lindsay Holmwood <lindsay@holmwood.id.au>"

# Stanzas
#
# Stanzas control when and how a process is started and stopped
# See a list of stanzas here: http://upstart.ubuntu.com/wiki/Stanzas#respawn

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Start the process
env GOMAXPROCS=16
exec /usr/local/bin/noodle /etc/noodle.conf
